# deploy/iso-installer/payload/master/kobeai-backup.timer
[Unit]
Description=Run KobeAI backup nightly at 02:00 local time
Requires=kobeai-backup.service

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true
RandomizedDelaySec=15min
Unit=kobeai-backup.service

[Install]
WantedBy=timers.target
